# Web/sql_based_cookie_injection/flask/flask.service
[Unit]
Description=flask
After=network.target

[Service]
User=flask
WorkingDirectory=/home/flask
ExecStart=/usr/local/bin/python waitress_server.py
Restart=always

[Install]
WantedBy=multi-user.target